# contrib/macos-launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.agentmux</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>agentmux</string>
    <string>serve</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/WORKDIR</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>AGENTMUX_HOST</key>
    <string>127.0.0.1</string>
    <key>AGENTMUX_PORT</key>
    <string>8765</string>
    <key>AGENTMUX_TOKEN</key>
    <string>REPLACE_WITH_A_RANDOM_TOKEN</string>
    <key>AGENTMUX_AGENT</key>
    <string>claude</string>
    <key>AGENTMUX_SESSION</key>
    <string>agentmux</string>
    <key>AGENTMUX_WINDOW</key>
    <string>repl</string>
    <key>AGENTMUX_WORKDIR</key>
    <string>/ABSOLUTE/PATH/TO/WORKDIR</string>
    <key>AGENTMUX_STARTUP_DELAY</key>
    <string>1.0</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/agentmux.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/agentmux.stderr.log</string>
</dict>
</plist>
